# config/supervisord/server.ini
[group:juhe]
programs = api-server, update-news
[program:api-server]
directory = /home/ubuntu/juhe_api
environment = PATH="/home/ubuntu/juhe_api/ENV/bin"
command = /home/ubuntu/juhe_api/ENV/bin/gunicorn app:api_app --config config/gunicorn_conf.py
process_name = api-server-%(process_num)s
autorestart = unexpected
redirect_stderr = true
stdout_logfile = /home/ubuntu/juhe_api/logs/api-server-stdout.log
stdout_logfile_maxbytes = 100MB
stdout_logfile_backups = 50
stdout_capture_maxbytes = 10MB
stderr_logfile = /home/ubuntu/juhe_api/logs/api-server-stderr.log
stderr_logfile_maxbytes = 100MB
stderr_logfile_backups = 50
stderr_capture_maxbytes = 10MB
startsecs = 5
stopwaitsecs = 30
stopasgroup = true